# gnome/autostart/Pipewire.desktop
[Desktop Entry]
Exec=pipewire &
Name=Pipewire
Type=Application
Version=1.0
